# cloud/templates/ansible/linux/roles/linuxprometheus/templates/nodeexporter.service.j2
[Unit]
Description=Prometheus Node Exporter

[Service]
Type=simple
ExecStart={{ nodeexporter_dir }}/node_exporter

[Install]
WantedBy=default.target
